# scripts/django_app.ini
[program:base]
command=/usr/local/bin/uwsgi /home/[username]/webapps/[proyect_name]/uwsgi.ini