# scripts/django_app.ini
[program:{{  APP_NAME }}]
command=/usr/local/bin/uwsgi /home/{{ username }}/webapps/{{ APP_NAME }}/uwsgi.ini